# deploy/aws/server/pp-indexer.service
[Unit]
Description=Private Payments Indexer
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/pp
# Contract IDs + start ledger are compiled defaults (fresh pool/ASP), but can
# be overridden here without a rebuild.
Environment=INDEXER_BIND=127.0.0.1:8080
Environment=DATABASE_URL=postgres://indexer:indexer@localhost:5434/indexer
Environment=INDEXER_RPC_URL=https://soroban-testnet.stellar.org
Environment=INDEXER_START_LEDGER=3408231
Environment=INDEXER_CONTRACTS=CDFXXZCDNFVQXMM6DUZWXHABGHCAWWXXA7T3IZYE7DY5ZLHKTR52VACV,CALBHI3CBBMEQ4CC57NA4FFEMM26TUHZPFXPJJ3DOJ4GKTK7BZA352BH,CB7MG5HSATWQ4S4C4TGWECCAIR66VM6M5KQEYOUCWUN6ZO5N3LBK4K7V
ExecStart=/opt/pp/bin/indexer
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
